# crontab
* * * * * python3 -u test.py >> /var/log/cron.log 2>&1
